# infra/app.service
[Unit]
Description=app

[Service]
Type=forking
ExecStart=/usr/bin/app

[Install]
WantedBy=multi-user.target